# src/crontab
*/3 * * * * python -u /opt/hakuna_matata/payments_app/src/hakuna_matata_payment_script.py > /dev/stdout
