# src/crontab
5 * * * * python /opt/hakuna_matata/payments_app/src/hakuna_matata_payment_script.py
